# labvm/scripts/files/etc/rc.local
#!/bin/sh
# SystemD executes rc.local at boot-time, if present
# So put here your system initialization commands

